# conf/vagrant/salt/roots/salt/test-smyt/supervisor.conf
[program:{{ app }}]
directory = {{ path }}
user = {{ user }}
command = {{ run_script }}
stdout_logfile = {{ log_file }}
stderr_logfile = {{ error_log_file }}