# conf/supervisor.conf
[program:test-smyt]
directory = /test_smyt
user = vagrant
command = /test_smyt/run.sh
stdout_logfile = /test_smyt/log/log.txt
stderr_logfile = /test_smyt/log/error_log.txt